# docker/keylime_agent/tpm_server.service
[Unit]
Description=TPM Server

[Service]
Type=simple
ExecStart=/usr/local/bin/tpm_server -rm
