# deployment/openshift-worker-nodes/tpm_server.service
[Unit]
Description=TPM Server

[Service]
Type=simple
WorkingDirectory=/tmp/
ExecStart=/etc/bin/tpm_server -rm

[Install]
WantedBy=multi-user.target

